# systemd-services/avr-manager.service
[Unit]
Description=AVR Input Switcher based on Bluetooth connections
After=multi-user.target

[Service]
Type=simple
User=pi
ExecStart=/usr/bin/python3 /home/pi/raspi-bt-hd-audio-receiver/avr-manager.py

[Install]
WantedBy=multi-user.target